# board/elebihan/demo-runit/fs-overlay/etc/sv/dropbear/run
#!/bin/sh
exec 2>&1
exec tcpsvd -v -i./peers 0 22 dropbear -i -E
