# board/elebihan/demo-runit/fs-overlay/etc/sv/getty-1/run
#!/bin/sh

exec getty -L tty1 115200 vt100
